# script/service/pretzel.target
[Unit]
Description=Pretzel Platform

[Install]
WantedBy=multi-user.target
