# script/service/pz-authd.service
[Unit]
Description=PZ Auth Daemon

After=pz-ipcd.service postgresql.service pz-db-ready.service
Requires=pz-ipcd.service
Wants=postgresql.service pz-db-ready.service
PartOf=pretzel.target

[Service]
Type=simple
ExecStart=/opt/pretzel/bin/pz-authd

Restart=always
RestartSec=3

LimitNOFILE=65536

StandardOutput=journal
StandardError=journal
SyslogIdentifier=pz-authd

[Install]
WantedBy=pretzel.target
