# script/service/pz-db-ready.service
[Unit]
Description=Pretzel DB Readiness Gate
# Block the pz daemons (which read their running-config from PostgreSQL) until the
# server is actually accepting connections. After= only orders against the unit being
# "active", not against the DB being ready — on a cold boot the server can be up as a
# unit but not yet listening, so the daemons would race it. This oneshot closes that
# gap. It always exits 0: if PG never comes up the daemons still start (degraded onto
# the startup-file fallback) and reconcile via Config::isStale() once PG is reachable.
After=postgresql.service
Wants=postgresql.service
PartOf=pretzel.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/bash -c 'for i in $(seq 1 60); do pg_isready -q -h 127.0.0.1 -p 5432 && exit 0; sleep 0.5; done; echo "pz-db-ready: PostgreSQL not ready after 30s — daemons will use file fallback"; exit 0'

[Install]
WantedBy=pretzel.target
