# script/service/pz-grafana.service
[Unit]
Description=Pretzel Grafana
After=pz-prometheus.service
Wants=pz-prometheus.service
PartOf=pretzel.target

[Service]
Type=simple
# grafana.ini reads admin_password via $__env{GF_SECURITY_ADMIN_PASSWORD}. Default is
# the localhost-dev "admin"; override here (or in a systemd drop-in) for production so
# the real password never lands in version control.
Environment=GF_SECURITY_ADMIN_PASSWORD=admin
ExecStart=/usr/share/grafana/bin/grafana server \
  --config=/etc/pretzel/grafana/grafana.ini \
  --homepath=/usr/share/grafana

StandardOutput=append:/var/log/pretzel/grafana.log
StandardError=append:/var/log/pretzel/grafana.log

Restart=always
RestartSec=2

[Install]
WantedBy=pretzel.target
