# script/service/pz-postgres-exporter.service
[Unit]
Description=Pretzel Postgres Exporter
After=network.target postgresql.service
Wants=postgresql.service
PartOf=pretzel.target

[Service]
Type=simple
# DATA_SOURCE_NAME embeds the DB password, so it is NOT hardcoded here. script/start.py
# generates /etc/pretzel/postgres-exporter.env (0600, gitignored) from PG_DB_PASSWORD —
# keeping it in lockstep with the role and out of version control. The pretzel role is
# granted pg_monitor at install time so the exporter can read the full pg_stat_* views.
EnvironmentFile=/etc/pretzel/postgres-exporter.env
ExecStart=/opt/pretzel/bin/postgres_exporter \
  --web.listen-address=0.0.0.0:9187

StandardOutput=append:/var/log/pretzel/postgres_exporter.log
StandardError=append:/var/log/pretzel/postgres_exporter.log

Restart=always
RestartSec=2

[Install]
WantedBy=pretzel.target
